# images/microvm/rootfs/etc/systemd/system/thaw-agent.service
[Unit]
Description=Firecracker Thaw Agent
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/thaw-agent
Restart=on-failure
RestartSec=5
User=root
Group=root

# Environment
Environment=MMDS_ENDPOINT=http://169.254.169.254
Environment=LOG_LEVEL=info
Environment=WORKSPACE_DIR=/workspace
Environment=RUNNER_DIR=/home/runner

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=thaw-agent

[Install]
WantedBy=multi-user.target

